# gitops/system/gitops-update.path
# VERSION=1
[Unit]
Description=Watch for GitOps trigger

[Path]
# Watch the file inside the relabeled directory
PathModified=%h/.config/containers/systemd/gitops/webhook/triggers/deploy/
Unit=gitops-update-scheduler.service

[Install]
WantedBy=default.target
