# gitops/system/gitops-update.service
[Unit]
Description=Run Update Script

[Service]
Type=oneshot
ExecStart=%h/.config/containers/systemd/gitops/scripts/gitops-update.sh
