# gitops/system/gitops-update.path
[Unit]
Description=Watch for GitOps trigger

[Path]
# Watch the file inside the relabeled directory
PathModified=%h/.config/containers/systemd/gitops/webhook/triggers/deploy/

[Install]
WantedBy=default.target
