# actual/actual.container
[Unit]
Description=Actual Budget - Personal Finance Manager
After=pocket-id.service

[Service]
EnvironmentFile=%h/.config/containers/systemd/.env
Restart=on-abnormal
TimeoutStartSec=60s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target default.target

[Container]
# Container properties
ContainerName=actual
Image=ghcr.io/actualbudget/actual:26.2.0
AutoUpdate=registry
Pull=newer

# Environment options
EnvironmentFile=.env

# Network options
Network=shared-network

# Volume binds
Volume=%h/container-data/actual:/data:Z

# Traefik Labels
Label="traefik.enable=true"
Label="traefik.http.routers.actual.rule=Host(`actual.${DOMAIN_URL}`)"
Label="traefik.http.routers.actual.entrypoints=web-secure,web-secure-internal"

# Glance Labels
Label="glance.name=Actual Budget"
Label="glance.icon=si:actualbudget"
Label="glance.url=https://actual.${DOMAIN_URL}"
Label="glance.description=Envelope Budgeting"
Label="glance.id=actual"
